# systemd/spamban-weekday-morning.service
[Unit]
Description=SpamBan check (weekday morning, 08:00–09:00 MSK)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=spamban
WorkingDirectory=/opt/spamban
EnvironmentFile=/opt/spamban/.env
Environment=SPAMBAN_RUN_LABEL=weekday-morning
Environment=SPAMBAN_WINDOW_MINUTES=60
ExecStart=/opt/spamban/venv/bin/python /opt/spamban/check.py
TimeoutStartSec=90min
Nice=10

[Install]
WantedBy=multi-user.target
